# utils/bootstrap-stage1/config/inittab
# -*- mode: sh -*-
#
# Boot-time system configuration/initialization script.
# This is run first except when booting in single-user mode.
#
# NOTE: Our rcS script mounts virtual filesystems (e.g., /dev, proc), creates
#       some device nodes, sets BusyBox's mdev as the hotplug helper, then
#       mounts the other filesystems in fstab.
::sysinit:/etc/init.d/rcS

# NOTE: The fist tty is special.  First of all, the device is not specified,
#       resulting in init autodetecting the console device (even works fine
#       with Xen's hvc0).  Instead of launching a getty or shell, however, it
#       runs ruckus_bootstrap_stage2's init script, followed by login.  The
#       stage2 script just immediately launches login if the special file
#       /.ruckus_bootstrap_stage2 isn't present.
#
# FIXME: might want /ruckus to be configurable somehow?
#
::respawn:-/ruckus/utils/bootstrap-stage2/init

# /sbin/getty invocations for selected ttys
#
# FIXME: Not sure what the baud rate should be, or if it even matters on a
#        non-serial tty.
#
# NOTE: These are completely pointless if we're booting up as a Xen domU... but
#       having them here doesn't seem to hurt anything and they are helpful if
#       you're actually booting up on real hardware.
tty2::respawn:/sbin/getty 38400 tty2
tty3::respawn:/sbin/getty 38400 tty3
tty4::respawn:/sbin/getty 38400 tty4
tty5::respawn:/sbin/getty 38400 tty5
tty6::respawn:/sbin/getty 38400 tty6

# Serial console (example only!)
#
# FIXME: This is messing with the Xen hvc0 console for some reason.  Only
#        enable this if you have a serial line and you're not using Xen.
#
#::respawn:/sbin/getty -L ttyS0 115200 vt100

# Stuff to do when restarting the init process
::restart:/sbin/init

# Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a
